# systemd/noisy-capture.service
[Unit]
Description=Noisy Pi Capture Daemon
Documentation=https://github.com/andjar/noisy_pi
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/bin/python3 /opt/noisy-pi/capture/capture_daemon.py
Restart=always
RestartSec=10

# Environment
Environment=NOISY_DATA_DIR=/var/lib/noisy-pi
Environment=NOISY_LOG_DIR=/var/log/noisy-pi
Environment=NOISY_CONFIG_DIR=/opt/noisy-pi/config

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening (but allow necessary access)
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/noisy-pi /var/log/noisy-pi /opt/noisy-pi/config
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target


