# systemd/noisy-web.service
[Unit]
Description=Noisy Pi Web Dashboard
Documentation=https://github.com/andjar/noisy_pi
After=network.target noisy-capture.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/noisy-pi/web
ExecStart=/usr/bin/php -S 0.0.0.0:8080 -t /opt/noisy-pi/web
Restart=always
RestartSec=5

# Environment
Environment=NOISY_DATA_DIR=/var/lib/noisy-pi
Environment=NOISY_CONFIG_DIR=/opt/noisy-pi/config

# Logging
StandardOutput=journal
StandardError=journal

# Security
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/noisy-pi /var/log/noisy-pi /opt/noisy-pi/config
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

